# cfg/supervisor_prod/conf.d/consumer_email.conf
[program:consumer_email]
command=php /srv/app/bin/console app:consumer --email --quiet
stdout_logfile=NONE
autostart=true
autorestart=true
startsecs = 0
startretries=1
stopsignal=KILL