# cfg/supervisor_prod/conf.d/consumer_package.conf
[program:consumer_package]
command=php /srv/app/bin/console app:consumer --package --quiet
stdout_logfile=NONE
autostart=true
autorestart=true
startsecs = 0
startretries=1
stopsignal=KILL